# web/supervisor/php-fpm.conf
[program:php-fpm]
command=php-fpm
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
